# docker/httpd-proxy/conf/cron/crontab
# Stats generator (every 5 mins)
*/5 * * * * /opt/bin/stats.sh
# keep this EOL